# scripts/com.lifey.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lifey.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/life/Documents/Codex/2026-06-19/refine-and-build-the-mvp-for/scripts/start-lifey.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>WorkingDirectory</key>
  <string>/Users/life/Documents/Codex/2026-06-19/refine-and-build-the-mvp-for</string>
  <key>StandardOutPath</key>
  <string>/Users/life/Library/Logs/Lifey/lifey.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/life/Library/Logs/Lifey/lifey-error.log</string>
</dict>
</plist>
